# zsh/profile.d/alias.plugin.zsh
# Custom plugin for zsh
#
# Common aliases
#
# Author: TimKing <timking.nf@foxmail.com>
# Date:   Sat, Jun 27  16:36:11 GMT+8 2020
#

###### Common Environments ######

# å¦‚æœæ²¡æœ‰æŒ‡å®šï¼Œåˆ™è‡ªåŠ¨é€‰æ‹©é¢œè‰²
LS_OPTIONS='--color=auto'
# æ˜¯å¦è¾“å‡ºé¢œè‰²
CLICOLOR='Yes'
# æŒ‡å®šé¢œè‰²
LSCOLORS='CxfxcxdxbxegedabagGxGx'


###### Command enhancements aliases ######

alias edit='nvim'
alias code='/Applications/Visual\ Studio\ Code.app/Contents/Resources/app/bin/code'

alias ls='ls -GFh'
alias ll="ls -l"
# åªæ˜¾ç¤ºéšè—æ–‡ä»¶
alias l.='ls -d .'
# å¤åˆ¶æ—¶ä¿ç•™é™å’Œä¿®æ”¹æ—¶é—´
alias cp='cp -v'
alias mv='mv -v'
alias mkdir='mkdir -pv'
alias less='less -FSRXc'
alias ssh='ssh -A'
alias ping='ping -c 5'
alias wget='wget -c'
# alias lc='colorls'
# alias lca='colorls -al'
alias scp="scp -p"


alias vper="${EDITOR:-edit} ~/.config/zsh/personal.zsh"
alias sper=". ~/.config/zsh/personal.zsh"
alias brash='~/.config/zsh/brash'
alias rm='brash -d -v'  # replace rm


###### Software ######

alias vim="nvim"
alias vi="nvim"


###### OSX alias ######

if [[ $CUSTOM_OS = "osx" ]]; then

alias date="gdate"  # use gnu date replace date
alias clearcp='printf "\e]1337;ClearCapturedOutput\e\\"'  # å‘Šè¯‰iterm2æ¸…ç† CapturedOutput
alias sed="gsed"  # use gnu sed replace osx sed

fi
